# f6-blanker.service
[Unit]
Description=Universal F6 Screen Blanker
# Ensures it only starts when the user's graphical session is ready
After=graphical-session.target

[Service]
# %h automatically becomes /home/username
ExecStart=/usr/bin/bash %h/.local/bin/f6-watcher.sh
Restart=always
# Inherits the user's environment automatically in --user mode
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
